# ops/zaocoworking-bot-autoredeploy.timer
[Unit]
Description=Poll origin/main every 10 min; redeploy bot if agent/ changed

[Timer]
OnBootSec=2min
OnUnitActiveSec=10min
Unit=zaocoworking-bot-autoredeploy.service
Persistent=true

[Install]
WantedBy=timers.target
